# plugui.ini
[program:plugui]
command=node server.js
directory=/opt/PlugUI
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/plugui.log
redirect_stderr=true